# conf.d/x11vnc.conf
[program:x11vnc]
command=/app/conf.d/x11vnc.sh
autorestart=true
